# systemd/dyndns.timer
[Unit]
Description=Periodically update Route53 A records to point here.

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
